# build/deb/stable/DEBIAN/postinst
#!/bin/sh

adduser --system --disabled-password tika

# start while booting
update-rc.d tika defaults

# start after installation
service tika restart
