# build/deb/stable/DEBIAN/prerm
#!/bin/sh

systemctl disable tika

systemctl stop tika

exit 0
