# DEBIAN/prerm
#!/bin/sh

systemctl stop tika

systemctl disable tika

exit 0
